# cookbooks/CLM/templates/default/JTS.conf.erb
description "JTS"

start on startup
start on runlevel [2345]
stop on runlevel [!2345]

env DISPLAY=localhost:1.0
export DISPLAY

post-start script
  exec Xvfb :1 -screen 0 800x600x24&
  exec su root -c "ulimit -n 65536;ulimit -u 10000;cd /opt/IBM/JazzTeamServer/server;./server.startup"
end script

post-stop script
  cd /opt/IBM/JazzTeamServer/server
  exec  ./server.shutdown
end script
