# cookbooks/CLM/templates/default/JTS.conf.erb
description "JTS"

start on startup
start on runlevel [2345]
stop on runlevel [!2345]

env DISPLAY=localhost:1.0

post-start script
	set -x
	Xvfb :1 -screen 0 800x600x24&
	exec su root -c "ulimit -n 65536;ulimit -u 10000;/opt/IBM/JazzTeamServer/server/server.startup"
end script

post-stop script
	set -x
	/opt/IBM/JazzTeamServer/server/server.shutdown
	pid=$(pidof /usr/bin/Xvfb)
	if [ -n "$pid" ]; then
		kill -9 $pid
	fi
end script
