# cookbooks/CLM/templates/default/JTS.conf.erb
description "JTS"

start on startup
start on runlevel [2345]

exec /opt/IBM/JazzTeamServer/server/server.startup